# ipipou.target
[Unit]
Description=IPIP-over-UDP tunnels via ipipou
